# omega2-base/files/etc/uci-defaults/20_onion_sysupgrade
#!/bin/sh

if ! grep ^'/root/' /etc/sysupgrade.conf; then
	echo "/root/" >> /etc/sysupgrade.conf
fi

exit 0
