# omega2-base/files/etc/uci-defaults/98_opkg
#!/bin/sh

. /etc/os-release

echo "src/gz omega2_onion http://repo.onioniot.com/omega2/packages/openwrt-${VERSION_ID}/onion" >> /etc/opkg/customfeeds.conf

exit 0
